# deploy/com.vy.inventory-backup.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.vy.inventory-backup</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/vy/Documents/Development/inventory-web/.venv/bin/python</string>
        <string>scripts/backup_db.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/vy/Documents/Development/inventory-web</string>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key><integer>3</integer>
        <key>Minute</key><integer>30</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/vy/Documents/Development/inventory-web/logs/backup.out.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/vy/Documents/Development/inventory-web/logs/backup.err.log</string>
</dict>
</plist>
