# deploy/com.vy.inventory.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.vy.inventory</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/vy/Documents/Development/inventory-web/.deploy/repo/.venv/bin/uvicorn</string>
        <string>app.main:app</string>
        <string>--host</string>
        <string>0.0.0.0</string>
        <string>--port</string>
        <string>8502</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/vy/Documents/Development/inventory-web/.deploy/repo</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>INVENTORY_DB</key>
        <string>/Users/vy/Documents/Development/inventory-web/inventory.db</string>
        <key>INVENTORY_ENV</key>
        <string>/Users/vy/Documents/Development/inventory-web/.env</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <!-- Slow the respawn loop when the port is taken or the venv is mid-rebuild,
         so a stuck service doesn't grow logs unbounded (learned the hard way). -->
    <key>ThrottleInterval</key>
    <integer>30</integer>

    <key>StandardOutPath</key>
    <string>/Users/vy/Documents/Development/inventory-web/logs/inventory.out.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/vy/Documents/Development/inventory-web/logs/inventory.err.log</string>
</dict>
</plist>
