# deploy/com.vy.inventory.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.vy.inventory</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/vy/Documents/Development/inventory-web/.venv/bin/uvicorn</string>
        <string>app.main:app</string>
        <string>--host</string>
        <string>0.0.0.0</string>
        <string>--port</string>
        <string>8502</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/vy/Documents/Development/inventory-web</string>

    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/vy/Documents/Development/inventory-web/logs/inventory.out.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/vy/Documents/Development/inventory-web/logs/inventory.err.log</string>
</dict>
</plist>
